# legacy/systemd/user/qwen-voice-input.service
[Unit]
Description=Qwen3-ASR push-to-talk voice input
After=graphical-session.target pipewire.service

[Service]
Type=simple
ExecStart=%h/AI/VoiceInput/run.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
